# aws/atmos-user-data/modules/atmos-user-data-helpers/templates/sync_iam_users.cron
#!/usr/bin/env bash

# fail fast
set -e

/opt/atmos/bin/sync_iam_users.sh
